# tools/exeq/tests/sd.redis.conf
# Note: redis should run in protected mode in a container if you are
# using it as a containerized service
[supervisord]
nodaemon=true
[program:redis]
; command=redis-server /etc/redis/redis.conf --protected-mode no --loglevel debug
command=/usr/local/bin/docker-entrypoint.sh
autostart=true
autorestart=false
stdout_logfile=/var/log/redis/stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/redis/stderr.log
stderr_logfile_maxbytes=10MB
exitcodes=0